# systemd/qzss-cloud-health-check.timer
[Unit]
Description=QZSS クラウド死活監視(30秒おき)

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
