# systemd/qzss-map.service
[Unit]
Description=QZSS 統合防災地図アプリ(server.js)
After=network.target
# クラッシュを繰り返しても再起動を止めない(既定では短時間に5回落ちると
# systemdが諦めて「failed」のまま放置してしまうため、上限を外す)。
# [Service]ではなく[Unit]セクションに書く必要がある
StartLimitIntervalSec=0

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/qzss/qzss-map
Environment=PORT=8080
# ラズパイ上のこのインスタンスはkiosk表示専用で、正本のデータはCloud Run
# 側にあるため、GCSへの読み書きを試みる必要が無い(サービスアカウント
# も無いので、試みても失敗ログが出るだけ)
Environment=LOCAL_STATE_ONLY=true
# 誰も公開サイトを見ていなくても、このラズパイ自身が計測したE2E
# レイテンシ(受信〜描画完了)を実績資料として貯めておきたいため、
# ローカルで完結した計測結果をクラウド側にも転送する
Environment=CLOUD_LATENCY_REPORT_URL=https://eq.shum10.com/client-timing
# ローカルkiosk表示ではトークン不要(ラズパイ内で完結するため)。
# Web配信も兼ねる場合はここに INGEST_TOKEN=xxxx を追記する
# 警報受信からキオスク画面への配信(WebSocket)を担う経路なので、
# バックグラウンドの定期タスクより優先させる
Nice=-5
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
