# systemd/qzss-kiosk.service
[Unit]
Description=QZSS キオスク表示(Chromium)
After=graphical.target qzss-map@%i.service
Wants=qzss-map@%i.service
# クラッシュを繰り返しても再起動を止めない([Service]ではなく[Unit]に書く)
StartLimitIntervalSec=0

[Service]
Type=simple
User=%i
# lightdmの自動ログインでXセッションが立ち上がるまで待つ。固定sleepだと
# 実際に早く準備できた場合でも常に同じ秒数待ってしまうため、ポーリング
# 方式にする(最大10秒まで)。ソケットの有無だけを見ると、Xサーバーが
# 起動した直後でまだ実際の解像度(EDID)ネゴシエーションが終わっていない
# タイミングでChromiumが起動し、意図しない低い解像度で表示されてしまう
# ことがある(実機で確認)。xrandrが実際の解像度を報告するまで待つことで
# これを避ける
ExecStartPre=/bin/bash -c 'export DISPLAY=:0; for i in $(seq 1 50); do xrandr --current 2>/dev/null | grep -qE " connected [0-9]" && exit 0; sleep 0.2; done'
# 警報描画がバックグラウンドの定期タスク(死活監視・OTAチェック等)と
# CPUを取り合わないよう、常にこちらを優先させる
Nice=-5
# Pi 3B+のGPUドライバ(VideoCore IV + Mesa)とChromium 150のGPU
# ラスタライズの組み合わせで、レンダラープロセスが数十分おきにクラッシュ
# することを実機で確認した(--type=rendererのクラッシュダンプ、GPU-GL
# コンテキスト関連)。ハードウェアGPUパスを完全に避け、ソフトウェア
# WebGL(SwiftShader)に固定して安定性を優先する。
# 注意: 単純な--disable-gpuだとWebGL自体が初期化に失敗し
# (MapLibre GL JSはWebGL必須、「Failed to initialize WebGL」で
# 地図が全く表示されなくなる)実機で確認済みのため使わないこと。
# --use-angle=swiftshader だけでは「Automatic fallback to software
# WebGL has been deprecated」で拒否されるため、
# --enable-unsafe-swiftshader との併用が必須(信頼済みの自前コンテンツ
# のみを表示するkioskなので、この安全性トレードオフは許容する)
# --disable-background-networking / --disable-sync / --disable-features=...:
# GoogleアカウントSync・翻訳・最適化ヒント等、kioskには不要な裏側の
# ネットワーク通信・機能を止める。実機のログでGCM(プッシュ通知)登録の
# 失敗("DEPRECATED_ENDPOINT"→"QUOTA_EXCEEDED")が繰り返し出ていたため、
# --disable-notifications(Push/Notification API自体を無効化)と
# --disable-features=GCM を追加した。ただしChromium側の既知の制約で
# GCM関連の内部初期化自体は完全には止まらない場合がある(CEFプロジェクトの
# 同種issueでも確立された完全な回避策は無いことを確認済み)。実害
# (クラッシュ等)は無いログの汚染・わずかな無駄通信程度のため、
# 完全に消えなくても許容する
# --js-flags="--max-old-space-size=128": Pi 3B+(実質855MB)はChromium
# だけで427〜742MBを使うことがあり常にメモリが逼迫している。V8のJS
# ヒープ上限を明示的に絞ることで、フレッシュ起動75秒後のRSSが実測で
# 742MB→655MB程度に下がることを確認した(ただし計測ごとの振れ幅が
# 大きいシステムなので目安程度)。なお同時に試した
# BackForwardCache/サイト分離無効化・Vulkan無効化・raster-threads削減
# の組み合わせは逆にメモリ使用量を増やしてしまい(427MB→660〜730MB)
# 却下した。Chromiumのフラグはこの機体では直感に反する結果になることが
# あるため、変更する際は1つずつ実測で検証すること
ExecStart=/bin/bash -c 'export DISPLAY=:0; export XDG_RUNTIME_DIR=/run/user/$(id -u %i); exec /usr/bin/chromium --kiosk --incognito --noerrdialogs --disable-infobars --disable-session-crashed-bubble --disable-restore-session-state --check-for-update-interval=31536000 --password-store=basic --use-angle=swiftshader --enable-unsafe-swiftshader --disable-smooth-scrolling --disable-background-networking --disable-sync --disable-notifications --disable-client-side-phishing-detection --disable-default-apps --js-flags="--max-old-space-size=128" --disable-features=Translate,MediaRouter,OptimizationHints,AutofillServerCommunication,GCM http://localhost:8080'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
